# deploy/oracle/daily-stock-briefing.timer
[Unit]
Description=Run Daily Stock Briefing at 22:30 KST

[Timer]
OnCalendar=*-*-* 22:30:00 Asia/Seoul
Persistent=true

[Install]
WantedBy=timers.target
