# deploy/oracle/daily-stock-briefing.timer
[Unit]
Description=Run Daily Stock Briefing at 08:00 KST

[Timer]
OnCalendar=*-*-* 08:00:00 Asia/Seoul
Persistent=true

[Install]
WantedBy=timers.target
